# etc/cron.daily/multistage-backup.sh
#!/bin/bash
set -e
/usr/sbin/all-backup.sh
/usr/sbin/backup-encrypt.sh
